# bootstrap/docker/postgres/supervisord.conf
[supervisord]
nodaemon = false

[program: postgres]
command = /etc/init.d/postgresql start